# cluster-agent/slurm-usage-history-exporter.timer
[Unit]
Description=SLURM Usage History Exporter Timer
Requires=slurm-usage-history-exporter.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00

# Run 15 minutes after boot (for systems that are down at scheduled time)
OnBootSec=15min

# If the system was off during scheduled time, run on next boot
Persistent=true

[Install]
WantedBy=timers.target
